# packaging/systemd/attune-notifier.service
[Unit]
Description=Attune WebSocket Notification Service
Documentation=https://git.rdrx.app/attune/attune
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=attune
Group=attune
ExecStart=/usr/bin/attune-notifier
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/attune/environment
WorkingDirectory=/var/lib/attune

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/attune /var/log/attune
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=attune-notifier

[Install]
WantedBy=multi-user.target
